# frontdeskagent/rootfs/etc/services.d/frontdeskagent/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the FrontDeskAgent service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting FrontDeskAgent config bridge"

INTEGRATION_SOURCE_DIR="/opt/frontdeskagent/custom_components/frontdeskagent"
INTEGRATION_TARGET_DIR="/config/custom_components/frontdeskagent"

if [ -d "${INTEGRATION_SOURCE_DIR}" ]; then
  mkdir -p "/config/custom_components"
  rm -rf "${INTEGRATION_TARGET_DIR}"
  cp -R "${INTEGRATION_SOURCE_DIR}" "${INTEGRATION_TARGET_DIR}"
  bashio::log.info "Synced FrontDeskAgent integration to ${INTEGRATION_TARGET_DIR}"
else
  bashio::log.warning "Integration source not found at ${INTEGRATION_SOURCE_DIR}"
fi

exec python3 /usr/bin/main.py
